# raspberry/server/cameraapi.service
[Unit]
Description=uWSGI instance to serve cameraapi
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/home/pi/remote-control-camera/raspberry/server
Environment="PATH=/home/pi/remote-control-camera/raspberry/server/cameraapiserver/bin"
ExecStart=/home/pi/remote-control-camera/raspberry/server/cameraapiserver/bin/uwsgi --ini cameraapi.ini

[Install]
WantedBy=multi-user.target
